# supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command="/usr/sbin/mysqld"

[program:grafana]
command="/usr/sbin/grafana-server"